# tt-bot.service
[Unit]
Description=TT Utilities Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/codex/tt-bot
ExecStart=/bin/bash /home/codex/tt-bot/start.sh
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
